# alpha-sniper/deployment/systemd/alpha-sniper-live-crash-notify.service
[Unit]
Description=Alpha Sniper LIVE Crash Notification
After=alpha-sniper-live.service

[Service]
Type=oneshot
User=alpha-sniper
Group=alpha-sniper
WorkingDirectory=/opt/alpha-sniper

# Send crash notification
ExecStart=/opt/alpha-sniper/venv/bin/python /opt/alpha-sniper/scripts/crash_notify.py

# Environment file for Telegram credentials
EnvironmentFile=/etc/alpha-sniper/alpha-sniper-live.env

[Install]
WantedBy=multi-user.target
